# configs/BSN.service
[Unit]
Description=BSN Bun Service
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=30s
Restart=on-failure
RestartSec=2s
TimeoutStartSec=30s
TimeoutStopSec=30s
WorkingDirectory=/opt/BSN/current
ExecStart=/usr/bin/bun run start
Environment=PORT=3000
Environment=NODE_ENV=production

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077

OOMScoreAdjust=-1000
LimitNOFILE=1048576
TasksMax=4096
CPUQuota=200%
MemoryMax=2G

KillSignal=SIGTERM
KillMode=mixed

[Install]
WantedBy=multi-user.target
